# rain_pipeline/crontab
00 00 * * * cd /rain_pipeline_cronjob && /usr/local/bin/python main.py >> cron.log 2>&1
